# supervisord.conf
[supervisord]
nodaemon=true

[program:bratconfig]
command=/bin/bash -c "/usr/bin/brat_install_wrapper.sh"
autorestart=false
startsecs=0


[program:apache2]
command=/bin/bash -c "source /etc/apache2/envvars && exec /usr/sbin/apache2 -DFOREGROUND"

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/sshd.log
stderr_logfile=/var/log/supervisor/sshd_err.log
